# pack/alphamon-agentpi_0.0-1/DEBIAN/alphamonpiagent.service
[Unit]
Description=ALPHAMONPIAGENT

[Service]
Type=simple
ExecStart=/usr/local/bin/http-client

[Install]
WantedBy=multi-user.target